# bin/com.fizzy.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fizzy.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/noviadihapsoro/Deployments/fizzy/bin/start-production</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/noviadihapsoro/Deployments/fizzy/log/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noviadihapsoro/Deployments/fizzy/log/launchd.error.log</string>
</dict>
</plist>
